# deploy/systemd/mcp-bidrento.service
[Unit]
Description=Bidrento MCP Server (@mcp-stack/bidrento)
After=network.target

[Service]
Type=simple
User=bdroapi
Group=bdroapi

# WorkingDirectory = package dir (where .env and .sessions.json live)
WorkingDirectory=/home/bdroapi/mcp-stack/packages/bidrento

# node resolves @mcp-stack/core via workspace symlinks from monorepo root
ExecStart=/usr/bin/node dist/index.js

Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/bdroapi/mcp-stack
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
